# nextcloud/nextcloud.cron
*/5 * * * *     docker exec --tty --user www-data flap_nextcloud php -f /var/www/html/cron.php
0   4 * * *     docker exec --tty --user www-data flap_nextcloud php occ preview:pre-generate
