# nextcloud/nextcloud.cron
*/5 * * * *     cd $FLAP_DIR; /usr/local/bin/docker-compose exec -T -u www-data nextcloud   php -f /var/www/html/cron.php
0   4 * * *     cd $FLAP_DIR; /usr/local/bin/docker-compose exec -T -u www-data nextcloud   php occ preview:pre-generate
 
