# system/system.cron
@midnight    manager update >> /var/log/flap/update.log

@hourly      manager healthcheck >> /var/log/flap/healthcheck.log

# Every minutes
* * * * *    manager handle_domain_request >> /var/log/flap/tls.log
