# system/system.cron
@reboot      flapctl setup cron

@midnight    flapctl update &> /var/log/flap/update.log

@hourly      flapctl healthcheck &> /var/log/flap/healthcheck.log
@hourly      flapctl setup network

# Every minutes
* * * * *    flapctl tls handle_request &> /var/log/flap/tls.log
* * * * *    flapctl tls update_dns_records &> /var/log/flap/tls.log
