# system/system.cron
@daily root /flap/system/scripts/update_system.sh
@daily root /flap/system/scripts/update_flap.sh
