# system/system.cron
@midnight  root    manager update system
@midnight  root    manager update flap

@hourly root    manager healthcheck

# Every minutes
* * * * *   root    manager tls applychanges
