# system/system.cron
@midnight    manager update flap

@hourly      manager healthcheck

# Every minutes
* * * * *    manager applychanges
