# system/system.cron
@reboot      flapctl ports setup > /dev/null

@midnight    flapctl update
@daily      flapctl update images

@hourly      flapctl ports setup > /dev/null

# Every minutes.
* * * * *    flapctl restart handle_request
* * * * *    flapctl domains handle_request
* * * * *    flapctl domains update_dns_records

# Every 30 minutes.
**/30 * * * *    flapctl users sync_mail_aliases

# Every day at 3 am and 15 pm.
0 3,15 * * *    flapctl backup
