# system/system.cron
@midnight    manager update &> /var/log/flap/update.log

@hourly      manager healthcheck &> /var/log/flap/healthcheck.log

# Every minutes
* * * * *    manager tls handle_request &> /var/log/flap/tls.log
* * * * *    manager tls update_dns_records &> /var/log/flap/tls.log
