# system/system.cron
@midnight    manager update flap > /var/log/flap/update.log

@hourly      manager healthcheck > /var/log/flap/healthcheck.log

# Every minutes
* * * * *    manager tls generate > /var/log/flap/tls.log
